# overlay/etc/systemd/system/ssh.service.d/generate-host-keys.conf
[Service]
ExecStartPre=/usr/bin/ssh-keygen -A
